# jenkins-slave-ruby/supervisord.conf
[supervisord]
nodaemon=true
[program:mysqld]
command=/usr/sbin/service mysql start
autorestart=false
[program:ssh]
commaand=/usr/sbin/service ssh start
autorestart=false
