# jenkins-slave-ruby/supervisord.conf
[program:mysqld]
command=/usr/sbin/service mysql start
autorestart=false